# configure/home/.config/systemd/user/files.timer
[Unit]
Description=Timer for checkfiles
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnBootSec=120
OnUnitActiveSec=600
Unit=files.service

[Install]
WantedBy=timers.target
