# configure/home/.config/systemd/user/files.service
[Unit]
Description=Run checkfiles program 

[Service]
Type=simple
ExecStart=%h/Programs/system/rofi/checkFiles.sh

[Install]
WantedBy=default.target
